# contrib/steel-clock.service
[Unit]
Description=steel-clock Nova Pro OLED daemon

[Service]
Type=simple
ExecStart=%h/.cargo/bin/steel-clock daemon --blank-on-exit
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
